# endpoint-hardening/10-detecting-host-compromise/data/alerts.json
[
  {
    "id": "alert-001",
    "timestamp": "2026-06-01T09:12:33Z",
    "agent": "meridian-payroll-01",
    "rule": {
      "id": "100001",
      "description": "SSH login from new IP address",
      "level": 5
    },
    "data": {
      "srcip": "10.20.30.45",
      "dstuser": "deploy",
      "program": "sshd",
      "action": "Accepted"
    }
  },
  {
    "id": "alert-002",
    "timestamp": "2026-06-01T09:13:01Z",
    "agent": "meridian-payroll-01",
    "rule": {
      "id": "100100",
      "description": "Cron modification detected — new crontab entry added",
      "level": 9
    },
    "data": {
      "audit": {
        "exe": "/usr/bin/crontab",
        "key": "cron_modification",
        "uid": "1001",
        "auid": "1001",
        "command": "crontab -e"
      },
      "file": "/var/spool/cron/crontabs/deploy",
      "entry": "*/5 * * * * /tmp/.update.sh"
    },
    "mitre": {
      "id": "T1053.003",
      "tactic": "Persistence",
      "technique": "Scheduled Task/Job: Cron"
    }
  },
  {
    "id": "alert-003",
    "timestamp": "2026-06-01T09:15:44Z",
    "agent": "meridian-payroll-01",
    "rule": {
      "id": "100200",
      "description": "Sensitive file read attempt: /etc/shadow",
      "level": 12
    },
    "data": {
      "audit": {
        "exe": "/usr/bin/python3",
        "key": "shadow_read",
        "uid": "1001",
        "auid": "1001",
        "syscall": "openat",
        "file": "/etc/shadow",
        "result": "EACCES"
      },
      "process": {
        "pid": 12345,
        "name": "python3",
        "cmdline": "python3 /tmp/.update.sh"
      }
    },
    "mitre": {
      "id": "T1003.008",
      "tactic": "Credential Access",
      "technique": "/etc/passwd and /etc/shadow"
    }
  },
  {
    "id": "alert-004",
    "timestamp": "2026-06-01T09:18:22Z",
    "agent": "meridian-payroll-01",
    "rule": {
      "id": "100300",
      "description": "SSH key added to authorized_keys for user deploy",
      "level": 10
    },
    "data": {
      "audit": {
        "exe": "/usr/bin/python3",
        "key": "ssh_key_added",
        "uid": "1001",
        "file": "/home/deploy/.ssh/authorized_keys",
        "operation": "write"
      },
      "process": {
        "pid": 12345,
        "name": "python3",
        "cmdline": "python3 /tmp/.update.sh"
      }
    },
    "mitre": {
      "id": "T1021.004",
      "tactic": "Lateral Movement",
      "technique": "Remote Services: SSH"
    }
  },
  {
    "id": "alert-005",
    "timestamp": "2026-06-01T09:20:01Z",
    "agent": "meridian-payroll-01",
    "rule": {
      "id": "100010",
      "description": "Routine package update check (benign)",
      "level": 3
    },
    "data": {
      "program": "apt",
      "action": "update",
      "packages": "ubuntu-advantage-tools"
    }
  },
  {
    "id": "alert-006",
    "timestamp": "2026-06-01T09:22:15Z",
    "agent": "meridian-payroll-01",
    "rule": {
      "id": "100011",
      "description": "Regular nginx access log entry (benign)",
      "level": 2
    },
    "data": {
      "srcip": "10.10.1.50",
      "url": "/api/v1/health",
      "status": 200,
      "program": "nginx"
    }
  }
]
