# ai-augmented-ops/07-ai-detection-triage/data/ground-truth.json
{
  "ALT-0001": {"severity": "CRITICAL", "technique": "T1059.001", "rationale": "Office spawning encoded PowerShell is a classic T1059.001 indicator; parent process makes this highly suspicious."},
  "ALT-0002": {"severity": "HIGH", "technique": "T1071.001", "rationale": "C2 over HTTPS to known Tor exit node; svchost as process is suspicious for user-initiated connections."},
  "ALT-0003": {"severity": "HIGH", "technique": "T1021.002", "rationale": "PsExec installation from a workstation to a file server indicates lateral movement."},
  "ALT-0004": {"severity": "MEDIUM", "technique": "T1110.001", "rationale": "Password spray pattern; 5 attempts is under lockout threshold — likely targeted, not automated spray."},
  "ALT-0005": {"severity": "MEDIUM", "technique": "T1052.001", "rationale": "Unregistered USB in Finance VLAN is a DLP concern; severity MEDIUM pending investigation of contents."},
  "ALT-0006": {"severity": "CRITICAL", "technique": "T1053.005", "rationale": "Scheduled task with IEX DownloadString is a textbook T1053.005 persistence + T1105 download; CRITICAL."},
  "ALT-0007": {"severity": "CRITICAL", "technique": "T1003.006", "rationale": "DCSync by non-replication account is credential dumping of the entire domain; maximum severity."},
  "ALT-0008": {"severity": "HIGH", "technique": "T1567.002", "rationale": "2.3 GB to personal Dropbox from Finance VLAN without DLP exception is likely data exfiltration."},
  "ALT-0009": {"severity": "CRITICAL", "technique": "T1490", "rationale": "Shadow copy deletion combined with unknown binary named ransomware.exe — active ransomware."},
  "ALT-0010": {"severity": "CRITICAL", "technique": "T1486", "rationale": "847 files encrypted with ransom extension in 30 seconds — active ransomware encryption event."},
  "ALT-0011": {"severity": "HIGH", "technique": "T1204.002", "rationale": "Emotet dropper quarantined but received 4 hours ago — need to determine if it executed first."},
  "ALT-0012": {"severity": "HIGH", "technique": "T1071.004", "rationale": "DNS to newly registered domain resolving to a suspicious IP is a C2 indicator; HIGH pending investigation."},
  "ALT-0013": {"severity": "HIGH", "technique": "T1021.004", "rationale": "SSH from non-VPN IP to jump host outside hours is suspicious; service account accessing jump host adds concern."},
  "ALT-0014": {"severity": "CRITICAL", "technique": "T1562.001", "rationale": "Disabling Windows Defender via registry by unknown binary from TEMP is defence evasion preceding further attack."},
  "ALT-0015": {"severity": "HIGH", "technique": "T1190", "rationale": "UNION SELECT patterns in SQL errors strongly indicate SQL injection attempt; HIGH for active exploitation."},
  "ALT-0016": {"severity": "CRITICAL", "technique": "T1003.001", "rationale": "Mimikatz signatures in LSASS memory is confirmed credential dumping; CRITICAL."},
  "ALT-0017": {"severity": "HIGH", "technique": "T1078", "rationale": "Impossible travel — two simultaneous sessions from different continents indicates account compromise."},
  "ALT-0018": {"severity": "CRITICAL", "technique": "T1105", "rationale": "Curl download and execute of unknown binary from non-Meridian IP is active malware staging."},
  "ALT-0019": {"severity": "CRITICAL", "technique": "T1136.001", "rationale": "Admin account created on DC by SYSTEM outside change window — post-compromise persistence on DC."},
  "ALT-0020": {"severity": "CRITICAL", "technique": "T1555.003", "rationale": "Browser credential store accessed by unknown process named credential_stealer.exe — active credential theft."},
  "ALT-0021": {"severity": "LOW", "technique": null, "rationale": "Legitimate signed Adobe update via AdobeARMservice. No security concerns."},
  "ALT-0022": {"severity": "LOW", "technique": null, "rationale": "User changing own password via self-service portal during business hours — normal activity."},
  "ALT-0023": {"severity": "LOW", "technique": null, "rationale": "Print spooler restart from SCCM with approved change ticket — normal IT activity."},
  "ALT-0024": {"severity": "LOW", "technique": null, "rationale": "Low disk space is an operational alert, not a security event."},
  "ALT-0025": {"severity": "LOW", "technique": null, "rationale": "Successful nightly backup — expected activity."},
  "ALT-0026": {"severity": "LOW", "technique": null, "rationale": "IT admin RDP from approved jump host — normal, recorded session."},
  "ALT-0027": {"severity": "HIGH", "technique": "T1547.001", "rationale": "VBScript in Startup folder created 3 hours ago is persistence; investigate origin before treating as LOW."},
  "ALT-0028": {"severity": "MEDIUM", "technique": "T1176", "rationale": "Unapproved browser extension with broad permissions — medium risk, may be adware or credential stealer."},
  "ALT-0029": {"severity": "HIGH", "technique": "T1140", "rationale": "certutil decode of executable from PowerShell is a known LOLBin technique for payload staging."},
  "ALT-0030": {"severity": "HIGH", "technique": "T1048.003", "rationale": "Bypassing internal mail relay to direct SMTP is an exfiltration-over-alternative-protocol indicator."},
  "ALT-0031": {"severity": "LOW", "technique": null, "rationale": "Planned deployment matching approved JIRA ticket — expected change."},
  "ALT-0032": {"severity": "CRITICAL", "technique": "T1055", "rationale": "Named pipe mimicking Windows system names by unknown process is process injection preparation."},
  "ALT-0033": {"severity": "CRITICAL", "technique": "T1134.001", "rationale": "SeImpersonatePrivilege token impersonation by non-system process from TEMP — privilege escalation."},
  "ALT-0034": {"severity": "HIGH", "technique": "T1562.004", "rationale": "Non-admin adding inbound RDP firewall rule is defence evasion / enabling lateral movement."},
  "ALT-0035": {"severity": "HIGH", "technique": "T1059.005", "rationale": "Excel macro executing Shell() to cmd.exe from email document is macro-based initial access."},
  "ALT-0036": {"severity": "CRITICAL", "technique": "T1048", "rationale": "Customer data exported via bcp to non-Meridian IP — active data exfiltration of sensitive data."},
  "ALT-0037": {"severity": "CRITICAL", "technique": "T1218.005", "rationale": "mshta.exe loading remote HTA from non-Meridian IP is LOLBin-based payload delivery."},
  "ALT-0038": {"severity": "LOW", "technique": null, "rationale": "Scheduled AV scan with no detections — expected operational activity."},
  "ALT-0039": {"severity": "HIGH", "technique": "T1185", "rationale": "Chrome remote debugging port opened by script allows browser session hijacking."},
  "ALT-0040": {"severity": "MEDIUM", "technique": null, "rationale": "Root container image is a misconfiguration; MEDIUM for policy violation, not active attack."},
  "ALT-0041": {"severity": "HIGH", "technique": "T1557.001", "rationale": "NTLMv1 when disabled by policy may indicate a downgrade attack (e.g. NTLM relay setup)."},
  "ALT-0042": {"severity": "LOW", "technique": null, "rationale": "Planned Patch Tuesday update via WSUS in approved window — expected."},
  "ALT-0043": {"severity": "CRITICAL", "technique": "T1003.001", "rationale": "LSASS memory dump by helpdesk account without approval is credential dumping."},
  "ALT-0044": {"severity": "HIGH", "technique": "T1112", "rationale": "Re-enabling WDigest stores plaintext passwords in memory — common pre-credential-dump step."},
  "ALT-0045": {"severity": "CRITICAL", "technique": "T1055.012", "rationale": "Process Hollowing into explorer.exe by unknown binary with no hash match — active code injection."},
  "ALT-0046": {"severity": "HIGH", "technique": "T1078", "rationale": "Slack login from unexpected city while user claims to be elsewhere — possible account compromise."},
  "ALT-0047": {"severity": "LOW", "technique": null, "rationale": "NTP drift below threshold — operational monitoring, not security event."},
  "ALT-0048": {"severity": "CRITICAL", "technique": "T1056.001", "rationale": "Luhn-valid card numbers pasted to external site in rapid succession — active card data exfiltration."},
  "ALT-0049": {"severity": "HIGH", "technique": "T1098", "rationale": "Privileged DB role granted without change ticket by DBA — possible insider threat or compromised DBA account."},
  "ALT-0050": {"severity": "CRITICAL", "technique": "T1562.001", "rationale": "EDR tamper — sensor stopped, real-time protection gone — classic pre-ransomware defensive evasion."}
}
